# static/flanneld.service.d/50-network-config.conf
[Unit]
Requires=etcd.service
[Service]
ExecStartPre=/usr/bin/etcdctl set /coreos.com/network/config  '{"Network":"4.0.0.0/8", "Backend": {"Type": "vxlan"}}'
